# config/systemd/usr/attendance-reader.service
[Unit]
Description=NFC Attendance Reader (PCSC)
After=pcscd.service
Wants=pcscd.service

[Service]
Type=simple
WorkingDirectory=%h/nfc
ExecStart=/usr/bin/python3 -u %h/nfc/core/attendance_reader.py
Restart=always
RestartSec=1
UMask=0022

[Install]
WantedBy=default.target
