# config/systemd/usr/attendance-payroll.timer
[Unit]
Description=attendance payroll daily

[Timer]
OnCalendar=*-*-* 15:05:00 UTC
Persistent=true

[Install]
WantedBy=timers.target
